# deploy-scripts/shiftpilot-deploy.timer
# ShiftPilot Auto-Deploy Timer
# Place in: /etc/systemd/system/shiftpilot-deploy.timer
# This will check for updates every 5 minutes

[Unit]
Description=ShiftPilot Auto-Deployment Timer
Requires=shiftpilot-deploy.service

[Timer]
# Run every 5 minutes
OnCalendar=*:0/5

# Run 1 minute after boot
OnBootSec=1min

# If the service fails, retry after 2 minutes
OnFailure=2min

# Ensure the timer is persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target
